# gbpservice/nfp/service_vendor_agents/haproxy/haproxy-agent/etc/init/haproxy_agent.conf
description "NVSD Haproxy on VM agent"
author "One Convergence Devices"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

chdir /var/run

exec start-stop-daemon --start --chuid root --exec /usr/bin/haproxy_agent -- -f /etc/haproxy_agent/haproxy_agent.cfg
