# gbpservice/nfp/bin/startup/redhat/nfp_orchestrator.service
[Unit]
Description=One Convergence NFP Orchestrator
After=syslog.target network.target

[Service]
User=neutron
ExecStart=/usr/bin/nfp  --config-file /etc/neutron/neutron.conf --config-file /etc/neutron/plugins/ml2/ml2_conf.ini --config-file /etc/neutron/nfp/nfp_orchestrator.ini --log-file /var/log/nfp/nfp_orchestrator.log

[Install]
WantedBy=multi-user.target

