# gbpservice/nfp/bin/startup/redhat/nfp_config_orch.service
[Unit]
Description=One Convergence NFP Config Orchestrator
After=syslog.target network.target

[Service]
Type=simple
User=neutron
ExecStart=/usr/bin/nfp --config-file /etc/neutron/nfp/nfp_config_orch.ini --config-file /etc/neutron/neutron.conf --log-file /var/log/nfp/nfp_config_orch.log

[Install]
WantedBy=multi-user.target

