# gbpservice/nfp/bin/startup/redhat/nfp_proxy.service
[Unit]
Description=One Convergence NFP Proxy
After=syslog.target network.target

[Service]
User=root
ExecStart=/usr/bin/nfpproxy_startup
Restart=on-abort

[Install]
WantedBy=multi-user.target

