# gbpservice/nfp/bin/startup/redhat/nfp_proxy_agent.service
[Unit]
Description=One Convergence NFP Proxy Agent
After=syslog.target network.target

[Service]
User=root
ExecStart=/usr/bin/nfp --config-file /etc/neutron/neutron.conf --config-file /etc/neutron/nfp/nfp_proxy_agent.ini --log-file /var/log/nfp/nfp_proxy_agent.log

[Install]
WantedBy=multi-user.target

