# gbpservice/nfp/bin/startup/redhat/nfp_pecan.service
[Unit]
Description=One Convergence NFP Rest Server
After=syslog.target network.target

[Service]
ExecStart=sudo pecan serve /usr/local/lib/python2.7/dist-packages/gbpservice/contrib/nfp/configurator/api/config.py
Restart=on-abort

[Install]
WantedBy=multi-user.target

